# ops/axiomforge-challenge.service
[Unit]
Description=AxiomForge Phase 8 grand challenge portfolio cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/axiomforge
Environment=AXIOMFORGE_ROOT=/root/axiomforge-state
Environment=AXIOMFORGE_PUBLICATION_REPO=/root/axiomforge-publications
Environment=AXIOMFORGE_PUBLICATION_BRANCH=autonomous-publications
ExecStart=/root/axiomforge/.venv/bin/axiomforge --root /root/axiomforge-state challenge-cycle --execute-route --timeout-seconds 30
