# ops/axiomforge-review.timer
[Unit]
Description=Run AxiomForge Phase 4 reviewer and replicator gate

[Timer]
OnBootSec=17min
OnUnitActiveSec=2h
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
