# ops/axiomforge.service
[Unit]
Description=AxiomForge autonomous research heartbeat
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/axiomforge
Environment=AXIOMFORGE_ROOT=/root/axiomforge-state
EnvironmentFile=-/root/axiomforge-state/secrets.env
ExecStart=/root/axiomforge/.venv/bin/axiomforge --root /root/axiomforge-state research-cycle --goal "maintain autonomous research memory"
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
