# ops/axiomforge-site.service
[Unit]
Description=AxiomForge Phase 5 public lab-note site builder
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/axiomforge
Environment=AXIOMFORGE_PUBLICATION_REPO=/root/axiomforge-publications
Environment=AXIOMFORGE_PUBLICATION_BRANCH=autonomous-publications
ExecStart=/root/axiomforge/.venv/bin/axiomforge site-build
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
